# scripts/config/supervisor.conf
[program:altexo-signal]
command=/var/www/altexo/signal/scripts/run
directory=/var/www/altexo/signal
user=altexo
autostart=true
autorestart=true

